# systemd/0pty-server.service
[Unit]
Description=0pty PTY persistence server
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
# Adjust this path if you install the binary somewhere else.
ExecStart=/usr/local/bin/0pty-server -b 127.0.0.1:6077 -- /bin/bash
Restart=always
RestartSec=1

[Install]
WantedBy=default.target
